# it8951-epaper.service
[Unit]
Description=IT8951 E-Paper Display API Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/it8951-epaper
ExecStart=/usr/bin/python3 /opt/it8951-epaper/server/server.py --host 0.0.0.0 --port 8888
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target